# UNICORN_FINAL/scripts/deepseek-loop.service
# DeepSeek Advisory Loop — systemd unit TEMPLATE
#
# This file is a TEMPLATE. It is NOT installed automatically by
# `.github/workflows/deploy.yml`. An operator must copy it to
# `/etc/systemd/system/deepseek-loop.service` on the host manually,
# review the EnvironmentFile, then `systemctl daemon-reload` +
# `systemctl enable --now deepseek-loop.service`.
#
# Safety:
#   * Runs as non-root user `unicorn`.
#   * Default-OFF: the script exits cleanly unless DEEPSEEK_LOOP_ENABLED=1
#     is set in the EnvironmentFile.
#   * Execute mode additionally requires DEEPSEEK_LOOP_EXECUTE=1 and
#     DEEPSEEK_LOOP_ADMIN_TOKEN. Without both it is advisory only (log).
#   * Restart=always with RestartSec=60 to recover from transient errors,
#     not to mask a misconfiguration.

[Unit]
Description=DeepSeek Advisory Loop (allowlist-only, default-off)
After=network.target unicorn-backend.service
Wants=network-online.target

[Service]
Type=simple
User=unicorn
Group=unicorn
EnvironmentFile=/opt/unicorn/.env
WorkingDirectory=/opt/unicorn
ExecStart=/usr/bin/node /opt/unicorn/scripts/deepseek-loop.js
Restart=always
RestartSec=60

# Hardening — defense in depth.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/unicorn/data/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Logging goes to journalctl + the per-script log file at
# /opt/unicorn/data/logs/deepseek-loop.log (override with DEEPSEEK_LOOP_LOG_PATH).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
